# config/systemd/user/battery-alert.service
[Unit]
Description=Notify battery status (service)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-alert
