# config/systemd/user/mailsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mailsync
